# rootfs/etc/cron.d/spamassassin
30 02 * * * root /usr/bin/sa-learn --ham /var/mail/vhosts/*/*/{{ VMAIL_SUBDIR }}/{{ FOLDER_HAM_LEARN_PATH }}/* | logger -t sa-learn-ham;
40 02 * * * root /usr/bin/sa-learn --spam /var/mail/vhosts/*/*/{{ VMAIL_SUBDIR }}/{{ FOLDER_SPAM_LEARN_PATH }}/* | logger -t sa-learn-spam
